# system/etc/udev/rules.d/99-power-management.rules
# Battery: Shift CPU to 'power' efficiency
SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="/usr/bin/sh -c 'echo power > /sys/devices/system/cpu/cpu*/cpufreq/energy_performance_preference'"

# Plugged In: Shift CPU to 'balance_performance'
SUBSYSTEM=="power_supply", ATTR{online}=="1", RUN+="/usr/bin/sh -c 'echo balance_performance > /sys/devices/system/cpu/cpu*/cpufreq/energy_performance_preference'"
